# bigclungus-meta/systemd/discord-ingest.timer
[Unit]
Description=Run Discord → Graphiti ingestion every 6 hours

[Timer]
# Run every 6 hours
OnCalendar=*-*-* 00,06,12,18:30:00
# Run once on boot if the timer was missed (e.g. during downtime)
Persistent=true
# Randomize start within 5 minutes to avoid thundering herd
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
