# bigclungus-meta/monitoring/alert-webhook.service
[Unit]
Description=Prometheus Alert -> BigClungus inject bridge
After=network.target

[Service]
ExecStart=/usr/bin/python3 /mnt/data/monitoring/alert_webhook.py
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=default.target
