# claude-config/systemd/user/integration-tests.service
[Unit]
Description=BigClungus Integration Tests

[Service]
Type=oneshot
ExecStart=/mnt/data/scripts/run-integration-tests.sh
