# claude-config/systemd/user/labs-diff-narrator.service
[Unit]
Description=Diff Narrator Lab (labs.clung.us/diff-narrator/)
After=network.target

[Service]
WorkingDirectory=/mnt/data/labs/diff-narrator
ExecStart=/home/clungus/.bun/bin/bun run src/index.ts
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
