# claude-config/systemd/user/memory-sweep.timer
[Unit]
Description=BigClungus Memory Sweeper timer — every 4 hours

[Timer]
OnCalendar=*-*-* 00,04,08,12,16,20:00:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
