# claude-config/systemd/user/temporal-worker.service
[Unit]
Description=Temporal Worker — Real Estate Listings
After=network.target temporal.service

[Service]
Type=simple
WorkingDirectory=/home/clungus/work/temporal-workflows
ExecStart=/usr/bin/python3 /home/clungus/work/temporal-workflows/worker.py
Restart=on-failure
RestartSec=10
ExecStopPost=/mnt/data/scripts/service-crash-alert.sh %n

# Load .env vars into the service environment
EnvironmentFile=/home/clungus/work/temporal-workflows/.env

StandardOutput=journal
StandardError=journal
SyslogIdentifier=temporal-worker

[Install]
WantedBy=default.target
